# ororo-etl.service
[Unit]
Description=Run the ororo-etl job

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run -it fetch-forecast
