# ororo-etl.timer
[Unit]
Description=Run ororo-etl 4 times a day (every 6 hours)

[Timer]
OnCalendar=*-*-* 0/6:00:00
Persistent=true

[Install]
WantedBy=timers.target
